# supervisord.addon.conf
; Include this file in the existing enshrouded-server supervisord config.
; It wires the addon RPC namespace and starts the gateway in the same container.

[rpcinterface:addon]
supervisor.rpcinterface_factory = supervisor_addon.rpcinterface:make_supervisor_rpcinterface
audit_log = /var/log/supervisor/supervisor-gateway-audit.log

[program:supervisor-gateway]
user=enshrouded
environment=HOME="/home/enshrouded",USER="enshrouded",LANG="en_US.UTF-8",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
command=supervisor-gateway
stdout_syslog=true
stderr_syslog=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=120
